# launchagents/com.claude.obsidian-sync-qdrant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.obsidian-sync-qdrant</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
        <string>/Users/YOUR_USERNAME/.claude/scripts/hooks/obsidian-sync-qdrant.py</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>OBSIDIAN_API_KEY</key>
        <string>YOUR_OBSIDIAN_API_KEY</string>
        <key>DASHSCOPE_API_KEY</key>
        <string>YOUR_DASHSCOPE_API_KEY</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/obsidian-sync-qdrant.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/obsidian-sync-qdrant.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
